# roles/setClusterCommon/files/rc.local
#!/bin/sh
iptables -P FORWARD ACCEPT
swapoff -a